# bbpi-paper.service
[Unit]
Description=BBPI paper-trade botu
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=utkua
WorkingDirectory=/home/utkua/bbpi
EnvironmentFile=-/home/utkua/bbpi/.env
ExecStart=/home/utkua/bbpi/venv/bin/python papertrade.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
